# debian/subutai-le.timer
[Unit]
Description=LetsEncrypt certificates updater timer

[Timer]
OnBootSec=1h
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
