# debian/subutai.service
[Unit]
Description=Subutai Daemon

[Service]
Type=simple
ExecStart=/usr/bin/subutai daemon
KillMode=process
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
Alias=subutai-agent.service
